# triai.target
[Unit]
Description=TriAI Multi-Agent Framework
After=network.target ollama.service
Wants=triai-messaging.service triai-agents.service ollama.service

[Install]
WantedBy=multi-user.target