# triai_deployment_20250812_102707/systemd/triai-messaging.service
[Unit]
Description=TriAI Messaging Server
After=network.target

[Service]
Type=simple
User=triai
Group=triai
WorkingDirectory=/opt/triai
ExecStart=/opt/triai/venv/bin/python main.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/triai

[Install]
WantedBy=multi-user.target
